# fedora-sway/etc/udev/rules.d/99-battery-brightness.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ATTR{online}=="0", RUN+="/usr/local/bin/auto-brightness.sh bat"
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ATTR{online}=="1", RUN+="/usr/local/bin/auto-brightness.sh ac"
